# src/gpumod/templates/systemd/vllm.service.j2
[Unit]
Description={{ service.name }}
After=network.target

[Service]
Type=simple
User={{ settings.user | default('root') }}
Environment="CUDA_VISIBLE_DEVICES={{ settings.cuda_devices | default('0') }}"
Environment="HF_HOME={{ settings.hf_home | default('/root/.cache/huggingface') }}"
{% for key, value in extra_env.items() %}
Environment="{{ key }}={{ value }}"
{% endfor %}

ExecStart={{ settings.vllm_bin | default('vllm') }} serve \
    {{ service.model_id }} \
    --port {{ service.port }} \
    --host 0.0.0.0 \
    --gpu-memory-utilization {{ unit_vars.gpu_mem_util | default(0.9) }} \
    --max-model-len {{ unit_vars.max_model_len | default(4096) }}{% if unit_vars.extra_args is defined %} \
    {{ unit_vars.extra_args }}{% endif %}

Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
